# bigdata/supervisord.d/server/alliedelec_recv_get_all_categories.ini
[program:alliedelec_recv_get_all_categories]
directory=/data/bigdata/server/recv_cron
command=python3.6 ./alliedelec/recv_get_all_categories.py
process_name=%(program_name)s
autostart=true
autorestart=true
startretries=3
user=root
redirect_stderr=true
stdout_logfile=/data/log/bigdata/alliedelec/recv_get_all_categories.log
stdout_logfile_maxbytes=100MB
