# bigdata/supervisord.d/server/avnet_recv_goods.ini
[program:avnet_recv_goods]
directory=/data/bigdata/server/recv_cron
command=python3.6 ./avnet/recv_goods.py
numprocs=1
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
startretries=3
user=root
redirect_stderr=true
stdout_logfile=/data/log/bigdata/avnet/recv_goods.log
stdout_logfile_maxbytes=100MB
