# bigdata/supervisord.d/client/alliedelec_recv_get_total_pages.ini
[program:alliedelec_recv_get_total_pages]
directory=/home/root/bigdata/client/recv_cron
command=python3.6 ./alliedelec/recv_get_total_pages.py
numprocs=1
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
startretries=3
user=root
redirect_stderr=true
stdout_logfile=/data/log/bigdata/alliedelec/%(program_name)s_%(process_num)02d.log
stdout_logfile_maxbytes=100MB
